# packaging/systemd/jukebox_ota_agent.service
[Unit]
Description=Jukebox OTA Agent — verificação de atualização
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jukebox-ota
Group=jukebox-ota
# ADR 0002: timer só check (HTTP + download + verify → ready_to_apply).
# Apply automático fica a cargo do kiosk após grace period; manual via systemd-run upgrade.
ExecStart=/opt/jukeeo/ota-agent/jukebox-ota-agent check --config /etc/jukeeo/ota-agent.json
# POC: exit 2 = update disponível (não é falha); exit 1 = erro real
SuccessExitStatus=0 2
SyslogIdentifier=jukebox-ota
StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jukebox/.local/share/com.jukeeo.kiosk/logs /home/jukebox/.local/share/com.jukeeo.kiosk /var/lib/jukebox-ota
StateDirectory=jukebox-ota
StateDirectoryMode=0750

[Install]
WantedBy=multi-user.target
