# packaging/systemd/jukebox_ota_agent.service
[Unit]
Description=Jukebox OTA Agent — verificação de atualização
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jukebox-ota
Group=jukebox-ota
ExecStart=/opt/jukeeo/ota-agent/jukebox-ota-agent check --config /etc/jukeeo/ota-agent.json
# POC: exit 2 = update disponível (não é falha); exit 1 = erro real
SuccessExitStatus=0 2
SyslogIdentifier=jukebox-ota
StandardOutput=journal
StandardError=journal

# Privilégio mínimo — ver docs/plans/PLANO_SEGURANCA_OTA_PENDENCIAS.md §3.4.1
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
StateDirectory=jukebox-ota

[Install]
WantedBy=multi-user.target
