# packaging/systemd/jukebox_ota_agent.service
[Unit]
Description=Jukebox OTA Agent — verificação e aplicação de atualização
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jukebox-ota
Group=jukebox-ota
# Respeita política SQLite (intervalo, janela, enabled); se houver update, faz download + apply.
ExecStart=/opt/jukeeo/ota-agent/jukebox-ota-agent upgrade --config /etc/jukeeo/ota-agent.json
SuccessExitStatus=0
SyslogIdentifier=jukebox-ota
StandardOutput=journal
StandardError=journal

# apply usa sudo -n systemctl (sudoers) — NoNewPrivileges impediria elevação.
NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
# Política OTA + estado partilhado com o kiosk; apply grava releases/backups e symlinks em /opt/jukeeo.
ReadWritePaths=/home/jukebox/.local/share/com.jukeeo.kiosk/logs /home/jukebox/.local/share/com.jukeeo.kiosk /opt/jukeeo/releases /opt/jukeeo/backups /opt/jukeeo
StateDirectory=jukebox-ota

[Install]
WantedBy=multi-user.target
